# dns/dnsRslv.init
#!/bin/sh /etc/rc.common
# https://openwrt.org/docs/techref/initscripts

START=99

USE_PROCD=1
PROG=/root/dnsRslv.lua

boot() {
 BOOT=1
 start "$@"
}

start_service() {
 procd_open_instance 1
 procd_set_param command $PROG
 procd_close_instance
}
